# jobs/com.vnstock.scheduler.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <!-- Job Label -->
    <key>Label</key>
    <string>com.vnstock.scheduler</string>

    <!-- Program to run -->
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/nghia.dinh/Projects/vietnam-stocks/jobs/start_scheduler.sh</string>
    </array>

    <!-- Working Directory -->
    <key>WorkingDirectory</key>
    <string>/Users/nghia.dinh/Projects/vietnam-stocks</string>

    <!-- Run on load -->
    <key>RunAtLoad</key>
    <true/>

    <!-- Keep alive (restart if crashes) -->
    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
    </dict>

    <!-- Standard output/error logs -->
    <key>StandardOutPath</key>
    <string>/tmp/vnstock_scheduler.out.log</string>

    <key>StandardErrorPath</key>
    <string>/tmp/vnstock_scheduler.err.log</string>

    <!-- Environment variables -->
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>

    <!-- Process type -->
    <key>ProcessType</key>
    <string>Background</string>

    <!-- Throttle interval (seconds to wait before restarting on crash) -->
    <key>ThrottleInterval</key>
    <integer>60</integer>
</dict>
</plist>
